# .config/autostart/startAudio.desktop
[Desktop Entry]
Type=Application
Name=T3RobotAudio
Exec=bash /home/pi/.config/autoStartAudio.sh
